# .config/systemd/user/daily-paper.timer
[Unit]
Description=Set your daily wallpaper
After=network-online.target

[Timer]
OnCalendar=daily
RandomizedDelaySec=1m
Persistent=true

[Install]
WantedBy=timers.target
